# provisioning/roles/neo4j/files/neo4j.conf
description "neo4j upstart"

start on startup
stop on shutdown

exec sudo /opt/neo4j/neo4j-community-{{ neo4j_version }}/bin/neo4j start
